# wazuh-bootstrap-api/deploy/systemd/wazuh-bootstrap-api.service
[Unit]
Description=Wazuh Bootstrap API
Documentation=file:/opt/wazuh-bootstrap-api/docs/DEPLOYMENT.md
After=network-online.target wazuh-manager.service
Wants=network-online.target

[Service]
Type=simple
User=wazuh-bootstrap
Group=wazuh-bootstrap
WorkingDirectory=/opt/wazuh-bootstrap-api
EnvironmentFile=/etc/wazuh-bootstrap-api.env
ExecStart=/opt/wazuh-bootstrap-api/.venv/bin/uvicorn app.main:app --host ${BIND_HOST} --port ${BIND_PORT} --workers ${UVICORN_WORKERS} --proxy-headers --forwarded-allow-ips ${FORWARDED_ALLOW_IPS}
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=20

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
InaccessiblePaths=-/var/ossec/etc/client.keys -/usr/share/wazuh-dashboard -/etc/shadow
UMask=0077

[Install]
WantedBy=multi-user.target
